# dsco/project_template/{{cookiecutter.project_name}}/builder/ansible/roles/jupyter-config/files/jupyter.conf
[program:jupyter]
command = jupyter notebook
user=root
directory = /srv
autostart = true
autorestart = true
stdout_logfile = /var/log/jupyter/jupyter.log
redirect_stderr = true
