# dsco/project_template/{{cookiecutter.project_name}}/builder/ansible/roles/data-refresh/files/data_refresh.conf
[program:data_refresh]
command = python update_data.py
user=root
directory = /srv
autostart = true
autorestart = true
stdout_logfile = /var/log/data_refresh.log
redirect_stderr = true
