# dsco/project_template/{{cookiecutter.project_name}}/builder/ansible/roles/uwsgi/files/uwsgi.conf
[program:uwsgi]
command = uwsgi --ini flask_config.ini
user=root
directory = /srv
autostart = true
autorestart = true
stdout_logfile = /var/log/uwsgi/flask_uwsgi.log
redirect_stderr = true
